# systemd-user/clawdbot-config-sync.path
[Unit]
Description=Watch Clawdbot config for changes

[Path]
PathChanged=%h/.clawdbot/clawdbot.json

[Install]
WantedBy=default.target
