# systemd-user/moltbot-config-sync-push.path
[Unit]
Description=Watch Moltbot config for changes (sync + push)

[Path]
PathChanged=%h/.clawdbot/clawdbot.json

[Install]
WantedBy=default.target
